# deploy/fa-insight-harvester.service
[Unit]
Description=FA Insight Harvester Web App
After=network.target postgresql.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/fa-insight-harvester
Environment=PATH=/home/%i/fa-insight-harvester/.venv/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/home/%i/fa-insight-harvester/.env
ExecStart=/home/%i/fa-insight-harvester/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
